# deploy/rust-tftp-server.service
[Unit]
Description=RustTFTP
Documentation=https://github.com/user/rust-tftp-server
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/rust-tftp --headless
ExecReload=/bin/kill -HUP $MAINPID
NotifyAccess=main
Restart=on-failure
RestartSec=5s
LimitNOFILE=4096

# Security hardening
User=tftp
Group=tftp
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/srv/tftp /var/log/rust-tftp.log
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=yes

# Capabilities: bind to port 69
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
